# curator4-1209152391.conf
[program:curator4-1209152391]
command=python3 /home/curator4/curator4-1209152391/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1209152391/test.err.log
stdout_logfile=/home/curator4/curator4-1209152391/test.out.log